# cron.d/check-composer-security
MAILTO=cpanel@alberon.co.uk
PATH=/bin:/usr/bin:/sbin

# Check for out of date Composer packages once a week
0 2 * * sun root /root/alberon-linux-toolkit/check-composer-security.sh

# vim: ft=crontab
