# cron.d/forward-system-mail
MAILTO=cpanel@alberon.co.uk

# Check all system accounts are forwarded to us once per hour
10 * * * *  root /root/alberon-linux-toolkit/forward-system-mail.sh

# vim: ft=crontab
