# cron.d/upgrade-composer
# Upgrade Composer at midnight
0 0 * * *  root /root/alberon-linux-toolkit/cron-upgrade-composer-run.sh

# vim ft=crontab
